# dist/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[program:murmur]
autorestart=true
command=murmurd -fg -ini %(ENV_MURMUR_CFG_FILE)s -supw %(ENV_MURMUR_SUPERUSER_PASSWORD)s
directory=%(ENV_MURMUR_HOME)s
stdout_logfile=syslog
stderr_logfile=syslog
user=%(ENV_MURMUR_USER)s

[program:cron]
autorestart=true
command=crond -f
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root


# vim: set filetype=dosini :
