# dist/supervisord.ini
[supervisord]
nodaemon=true

[program:murmur]
autorestart=true
command=murmur -fg -ini %(MURMUR_CFG_FILE)s
directory=%(DEFAULT_WORKDIR)s
stdout_logfile=syslog
stderr_logfile=syslog
user=%(MURMUR_USER)s

[program:cron]
autorestart=true
command=crond -f
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root

[program:rsyslog]
autorestart=true
command=rsyslogd -n
directory=/tmp
stdout_logfile=syslog
stderr_logfile=syslog
user=root
